# scripts/systemd/uce.service
[Unit]
Description=UCE FastCGI Runtime
After=network-online.target mariadb.service memcached.service uce.socket
Wants=network-online.target
Requires=uce.socket

[Service]
Type=simple
WorkingDirectory=/Code/uce.openfu.com/uce
RuntimeDirectory=uce
RuntimeDirectoryPreserve=yes
StateDirectory=uce
CacheDirectory=uce
ExecStartPre=/usr/bin/mkdir -p /var/cache/uce/work /var/lib/uce/uploads /var/lib/uce/sessions
ExecStart=/Code/uce.openfu.com/uce/bin/uce_fastcgi.linux.bin
ExecStartPost=/usr/bin/bash /Code/uce.openfu.com/uce/scripts/systemd/wait-ready.sh
Restart=always
RestartSec=2
TimeoutStopSec=30
KillMode=mixed
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
